# clouds/aws/global/quality-site/ansible/roles/github-runner/templates/github-runner.service.j2
[Unit]
Description=GitHub Actions Runner ({{ runner_name }})
After=network.target

[Service]
Type=simple
User={{ runner_user }}
Group={{ runner_group }}
WorkingDirectory={{ runner_home }}
ExecStart={{ runner_home }}/run.sh
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ runner_home }} /tmp

# Environment
Environment=RUNNER_ALLOW_RUNASROOT=false

[Install]
WantedBy=multi-user.target